# worker/obsworker-watch@.service
[Unit]
Description=obsworker statistics

[Service]
ExecStart=/usr/local/sbin/obsworker-watch
StandardInput=socket
