# worker/obsworker-watch.socket
[Unit]
Description=obsworker statistics

[Socket]
ListenStream=4242
Accept=Yes

[Install]
WantedBy=sockets.target
